# etc/systemd/system/showip.service
[Unit]
Description=Show IP
After=network.target
 
[Service]
Type=oneshot
ExecStart=/usr/bin/showip
 
[Install]
WantedBy=multi-user.target
